# systemd/locked-in.service
[Unit]
Description=Locked-In Daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Adjust the path below to your install location
ExecStart=%h/.local/bin/locked-in run-legacy
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
