# systemd/locked-in-web.service
[Unit]
Description=Locked-In Web Dashboard

[Service]
Type=simple
# Adjust the path below to your install location
ExecStart=%h/.local/bin/locked-in web --port 8765
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
